# ops/systemd/alphard.service
[Unit]
Description=Alphard trading cycle container
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=oneshot
EnvironmentFile=/etc/alphard/runner.env

# Keep below 15m so the next M15 slot has cleanup margin.
TimeoutStartSec=4min
TimeoutStopSec=30s
TimeoutStartFailureMode=terminate
KillMode=control-group
SendSIGKILL=yes
Restart=no

ExecStartPre=/bin/sh -c '/usr/bin/docker ps -aq --filter name=^/${ALPHARD_CONTAINER_NAME}$ | grep -q . && /usr/bin/docker rm -f ${ALPHARD_CONTAINER_NAME} || true'
ExecStart=/opt/alphard/bin/alphard-run-once.sh
ExecStopPost=/bin/sh -c '/usr/bin/docker ps -aq --filter name=^/${ALPHARD_CONTAINER_NAME}$ | grep -q . && /usr/bin/docker rm -f ${ALPHARD_CONTAINER_NAME} || true'
