# ops/systemd/alphard.service
[Unit]
Description=Alphard trading cycle container
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=oneshot
EnvironmentFile=/etc/alphard/runner.env

# Keep below 15m so the next M15 slot has cleanup margin.
TimeoutStartSec=14min
TimeoutStopSec=30s
TimeoutStartFailureMode=terminate
KillMode=control-group
SendSIGKILL=yes
Restart=no

# Recover from any orphaned previous container before starting a new slot.
ExecStartPre=-/usr/bin/docker rm -f ${ALPHARD_CONTAINER_NAME}

# No flock here.
ExecStart=/opt/alphard/bin/alphard-run-once.sh

# Always clean up after success, failure, or timeout.
ExecStopPost=-/usr/bin/docker rm -f ${ALPHARD_CONTAINER_NAME}