# ops/systemd/alphard-m30.timer
[Unit]
Description=Run Alphard at +1 minute of each M30 slot

[Timer]
OnCalendar=*-*-* *:01,31:00
AccuracySec=5s
Persistent=true
RandomizedDelaySec=0
Unit=alphard.service

[Install]
WantedBy=timers.target
